# ansible-roles/setup-nest-desktop/templates/nest-desktop.service.j2
[Unit]
Description=NEST Desktop
Requires=network.target
After=network.target

[Service]
ExecStart=/usr/bin/docker run --rm -dt -p 54286:54286 --name nest-desktop {{ nest_desktop_image }}

[Install]
WantedBy=network.target
